# src/deb/spotify-helios-master/prerm
#!/bin/sh

set -e

if [ -e "/etc/init/helios-master" ]; then
	invoke-rc.d helios-master stop || exit $?
fi
