# src/deb/spotify-helios-agent/prerm
#!/bin/sh

set -e

if [ -e "/etc/init/helios-agent.conf" ]; then
	invoke-rc.d helios-agent stop || exit $?
fi
